# scripts/lib/service-templates/moonraker.service
# Systemd service file for Moonraker
# Generated by gschpoozi

[Unit]
Description=Moonraker API Server
Documentation=https://moonraker.readthedocs.io/
Requires=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=%USER%
SupplementaryGroups=moonraker-admin
RemainAfterExit=yes
ExecStart=%HOME%/moonraker-env/bin/python \
    %HOME%/moonraker/moonraker/moonraker.py \
    -d %HOME%/printer_data
Restart=always
RestartSec=10
